# systemd/user/activitywatcher.service
[Unit]
Description=ActivityWatcher Service

[Service]
Type=simple
ExecStart=/opt/activitywatch/aw-qt
Restart=always
